# ops/daemons/worker/cape-worker.service
[Unit]
Description=CAPE Issue Worker Daemon (WORKER_ID)
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=USER_NAME
Group=USER_GROUP
WorkingDirectory=CAPE_ROOT_PATH/cape-worker

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/etc/default/cape-worker

# Main process
ExecStart=PYTHON_PATH WORKER_SCRIPT_PATH --worker-id WORKER_ID --poll-interval 10 --log-level INFO

# Restart policy
Restart=on-failure
RestartSec=30s

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=append:CAPE_ROOT_PATH/cape-worker/logs/worker_WORKER_ID_stdout.log
StandardError=append:CAPE_ROOT_PATH/cape-worker/logs/worker_WORKER_ID_stderr.log
SyslogIdentifier=cape-worker-WORKER_ID

# Process priority
Nice=5

[Install]
WantedBy=multi-user.target
