# ops/daemons/worker/com.cape.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service Label -->
    <key>Label</key>
    <string>com.cape.worker.WORKER_ID</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>PYTHON_PATH</string>
        <string>WORKER_SCRIPT_PATH</string>
        <string>--worker-id</string>
        <string>WORKER_ID</string>
        <string>--poll-interval</string>
        <string>10</string>
        <string>--log-level</string>
        <string>INFO</string>
    </array>

    <!-- Working Directory -->
    <key>WorkingDirectory</key>
    <string>CAPE_ROOT_PATH/cape-worker</string>

    <!-- Environment Variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Run at load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive - restart if crashed -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Standard output and error logs -->
    <key>StandardOutPath</key>
    <string>CAPE_ROOT_PATH/cape-worker/logs/worker_WORKER_ID_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>CAPE_ROOT_PATH/cape-worker/logs/worker_WORKER_ID_stderr.log</string>

    <!-- Throttle restart interval (seconds) -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Process type -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
